# com.shapeandship.podcast-ripper.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shapeandship.podcast-ripper</string>
    <key>ProgramArguments</key>
    <array>
        <!-- caffeinate (-i idle, -s system) keeps the Mac awake for the whole
             run and releases the moment rip.py exits. Prevents sleep from
             freezing a long transcription mid-batch. -->
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>-s</string>
        <string>__PYTHON_BIN__</string>
        <string>__PROJECT_DIR__/rip.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/podcast-ripper.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/podcast-ripper.log</string>
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
